# deploy/systemd/turnstone.target
[Unit]
Description=Turnstone services (server + console)
# Group-target so operators can `systemctl start turnstone.target` and
# bring the required stack up/down. Individual units list
# `PartOf=turnstone.target` so a stop on the target propagates to them.
#
# turnstone-channel.service is deliberately NOT pulled in here — it's
# optional, and a `Wants=` reference here would resurrect it on every
# `systemctl start turnstone.target` even after `disable --now`. The
# channel unit has its own `WantedBy=multi-user.target` install line,
# so an `enable` makes it autostart at boot independently.
Requires=turnstone-server.service turnstone-console.service
After=turnstone-server.service turnstone-console.service

[Install]
WantedBy=multi-user.target
